# scripts/systemd/sparkle-tunnel.service
[Unit]
Description=Sparkle Cloudflare Tunnel
After=network.target sparkle.service

[Service]
Type=simple
User=YOUR_USER
ExecStart=/home/YOUR_USER/.local/bin/cloudflared tunnel run sparkle
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
