# scripts/systemd/sparkle-tunnel.service
[Unit]
Description=Sparkle Cloudflare Tunnel
After=network.target sparkle.service
StartLimitBurst=5
StartLimitIntervalSec=120

[Service]
Type=simple
# TODO: Change to your Linux username
User=YOUR_USER
MemoryMax=256M
# cloudflared path — auto-detected by install-services.sh
ExecStart=CLOUDFLARED_BIN --config /home/YOUR_USER/.cloudflared/sparkle-config.yml tunnel run sparkle
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
